# etc/cron.d/calendar
0,30 * * * *      root /usr/bin/php /var/www/cgi-bin/cacheCalendarJson.php
